# LSST-Alert-QA/systemd/sso-monitor.timer.example
# Rename to sso-monitor.timer

[Unit]
Description=Run SSO monitor daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target